# packages/yabu/yabu-0.1.3-py3-none-any.whl/yabu-0.1.3.data/data/yabu.timer
[Unit]
Requires=yabu.service

[Timer]
Unit=yabu.service
OnCalendar=weekly

[Install]
WantedBy=timers.target