# packages/yabu/yabu-0.1.2-py3-none-any.whl/yabu-0.1.2.data/data/yabu.service
[Unit]
Description=An utility to automate backups
After=network.targets

[Service]
ExecStart=/usr/bin/python3 -m yabu
Type=oneshot

[Install]
WantedBy=multi-user.target