# packages/tessdb-server/tessdb-server-2.6.8.tar.gz/tessdb-server-2.6.8/files/lib/systemd/system/tessdb.service
[Unit]
Description=TESS Quality Sky Meter database service
Documentation=https://github.com/STARS4ALL/tessdb-server

[Service]
Type=simple
User=root
KillMode=process
ExecStart=/usr/bin/python -m tessdb --config /etc/tessdb/config --log-file /var/log/tessdb.log
ExecReload=/bin/kill -s HUP -- $MAINPID

[Install]
WantedBy=multi-user.target
