# packages/genocide/genocide-9.tar.gz/genocide-9/files/genocide.service
[Unit]
Description=GENOCIDE - the king of the netherlands commits genocide OTP-CR-117/19/001 otp.informationdesk@icc-cpi.int https://genocide.rtfd.io

After=network-online.target
Wants=network-online.target

[Service]
User=genocide
Group=genocide
ExecStart=/usr/local/bin/genocided mods=irc,rss -b

[Install]
WantedBy=multi-user.target
