# packages/kpsh/kpsh-0.5.0.tar.gz/kpsh-0.5.0/contrib/kpsh.service
[Unit]
Description=KeePass Shell Daemon

[Service]
ExecStart=%h/.local/bin/kpsh -d
Restart=always

[Install]
WantedBy=default.target
